# ansible/roles/app_deploy/files/load-secrets.service
[Unit]
Description=Load secrets from AWS Secrets Manager
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/load-secrets.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
